# templates/net/house_net.service.j2
[Unit]
After=systemd-networkd.service
#After=docker.service

[Service]
Type=oneshot
# haven't yet debugged why this doesn't work right away
ExecStart=bash -c "sleep 10; /usr/sbin/iptables -A POSTROUTING --table nat --out-interface {{out_interface}} --jump MASQUERADE"

[Install]
WantedBy=multi-user.target
